# n8n-daemon.service
[Unit]
Description=n8n Daemon Service
Documentation=https://github.com/eskemo-hub/n8n-instance-manager
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/n8n-daemon/daemon
ExecStart=/usr/bin/node /opt/n8n-daemon/daemon/dist/index.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=n8n-daemon

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/opt/n8n-daemon /var/run/docker.sock /run/docker.sock
SupplementaryGroups=docker

# Environment
Environment=NODE_ENV=production

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target
